# config/launchd/com.aetherclaude.tinyproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aetherclaude.tinyproxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/tinyproxy</string>
        <string>-d</string>
        <string>-c</string>
        <string>/opt/homebrew/etc/tinyproxy/tinyproxy.conf</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/aetherclaude/logs/tinyproxy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aetherclaude/logs/tinyproxy-error.log</string>
</dict>
</plist>
